# docker/php/supervisor/main.conf
[program:php]
command = php-fpm
numprocs=1
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d